# deploy/spy-der-shadow.service
[Unit]
Description=SPY-DER parallel-track shadow publisher (paper/shadow only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zerodte
Group=zerodte
WorkingDirectory=/opt/spy-der
EnvironmentFile=-/etc/zerodte/zerodte.env
# systemd creates /var/lib/zerodte owned by the service user before start, so
# the runner can write the state file without root-owned /var/lib blocking the
# first-run mkdir.
StateDirectory=zerodte
ExecStart=/opt/zerodte/venv/bin/spy-der vps-runner \
    --live-state %S/zerodte/spy_der_state.json \
    --interval 60
Restart=always
RestartSec=5

# Live trading is forbidden. This unit only publishes observation state.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
